# autostart/mpnotd.desktop
[Desktop Entry]
Name=MPD Notification Daemon
Type=Application
Exec=mpnotd
TryExec=mpnotd
StartupNotify=false
