# autostart/mpnotd.desktop
[Desktop Entry]
Name=mpnotd
Type=Application
Exec=mpnotd
StartupNotify=false
